# supervisor_services.conf
[program:anope]
user=ircd
priority=60
autostart=true
autorestart=true
command=/home/ircd/unrealircd/services/bin/services -nofork 

[group:irc]
programs=anope
